# luci-app-mikhmonv3/root/etc/uci-defaults/99-mikhmonv3
#!/bin/sh

if command -v php8 &> /dev/null; then
    PHPVER="php8"
elif command -v php7 &> /dev/null; then
    PHPVER="php7"
else
    echo "PHP tidak ditemukan! Harap install php7 atau php8."
    exit 1
fi

opkg update
opkg install ${PHPVER}-cli ${PHPVER}-mod-xml ${PHPVER}-mod-filter ${PHPVER}-mod-curl

if [ -e /www/mikhmonv3 ]; then
    rm -rf /www/mikhmonv3
fi

ln -s /usr/share/mikhmonv3 /www/mikhmonv3

echo "Konfigurasi selesai. Menggunakan ${PHPVER}."
