# luci-app-mikhmonv3/root/etc/uci-defaults/99-mikhmonv3
#!/bin/sh

if command -v php &>/dev/null; then
    PHPVER=$(php -v | grep -oE '^PHP [0-9]+' | awk '{print "php" $2}')
elif command -v php8 &>/dev/null; then
    PHPVER="php8"
elif command -v php7 &>/dev/null; then
    PHPVER="php7"
else
    PHPVER=$(opkg list-installed | grep -oE '^php[0-9]' | head -n1)
fi

if [ -z "$PHPVER" ]; then
    echo "❌ Tidak dapat mendeteksi PHP. Pastikan php7 atau php8 sudah diinstal."
    exit 1
fi

echo "✅ PHP terdeteksi: $PHPVER"

opkg update
opkg install ${PHPVER}-cli ${PHPVER}-mod-xml ${PHPVER}-mod-filter ${PHPVER}-mod-curl

[ -e /www/mikhmonv3 ] && rm -rf /www/mikhmonv3
ln -s /usr/share/mikhmonv3 /www/mikhmonv3
echo "✅ Mikhmon v3 siap diakses di /www/mikhmonv3"
