# luci-app-mikhmonv3/root/etc/uci-defaults/99-mikhmonv3
#!/bin/sh

[ -e /www/mikhmonv3 ] && rm -rf /www/mikhmonv3; ln -s /usr/share/mikhmonv3 /www/mikhmonv3